# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/streambox/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

# Virtual display - must start first
[program:xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_DISPLAY_WIDTH)sx%(ENV_DISPLAY_HEIGHT)sx%(ENV_DISPLAY_DEPTH)s -ac +extension GLX +render -noreset
autorestart=true
priority=10
stdout_logfile=/var/log/streambox/xvfb.log
stderr_logfile=/var/log/streambox/xvfb-err.log

# PulseAudio virtual audio
[program:pulseaudio]
command=/usr/bin/pulseaudio --daemonize=no --exit-idle-time=-1 -F /tmp/pulse/default.pa
autorestart=true
priority=15
stdout_logfile=/var/log/streambox/pulse.log
stderr_logfile=/var/log/streambox/pulse-err.log

# Window manager
[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY="%(ENV_DISPLAY)s"
autorestart=true
priority=20
startsecs=3
stdout_logfile=/var/log/streambox/fluxbox.log
stderr_logfile=/var/log/streambox/fluxbox-err.log

# VNC server
[program:x11vnc]
command=/usr/bin/x11vnc -display %(ENV_DISPLAY)s -rfbauth /root/.vnc/passwd -rfbport 5900 -shared -forever -noxdamage -noxfixes -noxrecord
autorestart=true
priority=25
startsecs=3
stdout_logfile=/var/log/streambox/x11vnc.log
stderr_logfile=/var/log/streambox/x11vnc-err.log

# noVNC web client
[program:novnc]
command=/opt/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autorestart=true
priority=30
startsecs=3
stdout_logfile=/var/log/streambox/novnc.log
stderr_logfile=/var/log/streambox/novnc-err.log

# SSH server
[program:sshd]
command=/usr/sbin/sshd -D -p 2222
autorestart=true
priority=30
stdout_logfile=/var/log/streambox/sshd.log
stderr_logfile=/var/log/streambox/sshd-err.log

# FFmpeg RTMP stream - auto-starts with container
[program:ffmpeg-stream]
command=/opt/streambox/start-stream.sh
autorestart=false
startretries=3
priority=40
startsecs=0
stdout_logfile=/var/log/streambox/stream-supervisor.log
stderr_logfile=/var/log/streambox/stream-supervisor-err.log

# Stream control panel GUI
[program:stream-panel]
command=/usr/bin/python3 /opt/streambox/stream-panel.py
environment=DISPLAY="%(ENV_DISPLAY)s",RTMP_URL="%(ENV_RTMP_URL)s",STREAM_KEY="%(ENV_STREAM_KEY)s"
autorestart=true
priority=50
startsecs=5
stdout_logfile=/var/log/streambox/panel.log
stderr_logfile=/var/log/streambox/panel-err.log
